# scripts/launchd/com.jigger.everythingsearch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jigger.everythingsearch</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <!-- 安装脚本 install.sh 会自动生成此路径，手动安装请替换为你的实际路径 -->
        <string>/Users/YOUR_USERNAME/.local/bin/everythingsearch_index.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>1800</integer>
</dict>
</plist>
